# 3rdparty/libftd2xx/libftd2xx_linux/chore/85-unbind-ftdi.rules
ACTION=="add", ATTRS{idVendor}=="0403", MODE:="666"
ACTION=="add", ATTRS{idVendor}=="0403", ATTR{interface}=="Single RS232-HS", RUN+="/bin/sh -c '/bin/echo $kernel > /sys/bus/usb/drivers/ftdi_sio/unbind'"
